# roles/prepare/templates/xroad-scanner.service
[Unit]
Description=X-Road Network Scanner Service

[Service]
Type=simple
RemainAfterExit=no
ExecStart=/usr/bin/python3 /tmp/server.py
